# Chapter01/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:myserver]
command=/home/sod/go/bin/romanserver
autostart=true
autorestart=true
redirect_stderr=true
